# eclipse-che-installer-osx/src/main/app/Library/LaunchDaemons/com.codenvy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--

    Copyright (c) 2015 Codenvy, S.A.
    All rights reserved. This program and the accompanying materials
    are made available under the terms of the Eclipse Public License v1.0
    which accompanies this distribution, and is available at
    http://www.eclipse.org/legal/epl-v10.html

    Contributors:
      Codenvy, S.A. - initial API and implementation

-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>UserName</key>
    <string>root</string>

    <key>Label</key>
    <string>com.codenvy</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Applications/EclipseChe${eclipse-che.version}.app/Contents/Resources/eclipse-che/bin/che.sh</string>
      <string>run</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

  </dict>

</plist>
